# ops/systemd/espo-backup-audit-prod.timer
[Unit]
Description=Daily EspoCRM prod backup audit

[Timer]
Unit=espo-backup-audit@prod.service
OnCalendar=*-*-* 03:50:00
Persistent=true
AccuracySec=1m

[Install]
WantedBy=timers.target
