# ops/systemd/espo-backup-audit-dev.timer
[Unit]
Description=Daily EspoCRM dev backup audit

[Timer]
Unit=espo-backup-audit@dev.service
OnCalendar=*-*-* 04:20:00
Persistent=true
AccuracySec=1m

[Install]
WantedBy=timers.target
